# crontab
* * * * * /bin/sh -c "date >> /var/log/cron.log 2>&1"
